# deploy/com.agentharness.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agentharness</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>-m</string>
        <string>agentharness.cli</string>
        <string>serve</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>AGENTHARNESS_HOME</key>
        <string>__HARNESS_HOME__</string>
        <key>PATH</key>
        <string>__PATH__</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>__HOME__</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HARNESS_HOME__/logs/serve.out.log</string>

    <key>StandardErrorPath</key>
    <string>__HARNESS_HOME__/logs/serve.err.log</string>
</dict>
</plist>
